# ansible/roles/vpn/templates/supervisor/vpn.conf
[program:vpn-server]
command = /opt/vpn/libexec/server
autostart = false
redirect_stderr = true
